# roles/edpm_neutron_dhcp/templates/quadlet/neutron_dhcp_agent.container.j2
# Quadlet container file for neutron_dhcp_agent
# Managed by edpm_container_standalone - do not edit manually
# Requires Podman 4.6+ (RHEL 9.4+)
{%- set _dhcp_volumes = [] -%}
{%- set _dhcp_volumes = _dhcp_volumes
    + edpm_neutron_dhcp_common_volumes
    + edpm_neutron_dhcp_tls_cacert_volumes %}


# DO NOT REMOVE - used by container_config_hash module
# config_hash=

[Unit]
Description=Neutron DHCP Agent container (EDPM)

[Container]
ContainerName=neutron_dhcp_agent
Image={{ edpm_neutron_dhcp_image }}
Network=host
User=root
PodmanArgs=--privileged
PodmanArgs=--pid=host
PodmanArgs=--cgroupns=host
LogDriver=journald
Environment=KOLLA_CONFIG_STRATEGY=COPY_ALWAYS
{% for vol in _dhcp_volumes -%}
Volume={{ vol }}
{% endfor -%}
Label=managed_by=edpm_ansible

[Service]
Restart=always
TimeoutStartSec=900
TimeoutStopSec=84

[Install]
WantedBy=multi-user.target
