# roles/edpm_tripleo_cleanup/molecule/default/test-data/fake-tripleo-service.service
[Unit]
Description=Fake tripleo service

[Service]
ExecStart=/usr/bin/true

[Install]
WantedBy=multi-user.target
