# roles/edpm_telemetry/templates/quadlet/ceilometer_agent_compute.container.j2
# Quadlet container file for ceilometer_agent_compute
# Managed by edpm_container_standalone - do not edit manually
# Requires Podman 4.6+ (RHEL 9.4+)
{%- set _ceilometer_volumes = edpm_telemetry_ceilometer_common_volumes -%}
{%- if ca_bundle_exists | bool -%}
{%- set _ceilometer_volumes = _ceilometer_volumes +
        edpm_telemetry_ceilometer_tls_cacert_volumes -%}
{%- endif -%}
{%- if tls_cert_exists | bool -%}
{%- set _ceilometer_volumes = _ceilometer_volumes +
        edpm_telemetry_ceilometer_tls_volumes -%}
{%- endif %}

# DO NOT REMOVE - used by container_config_hash module
# config_hash=

[Unit]
Description=Ceilometer Agent Compute container (EDPM)

[Container]
ContainerName=ceilometer_agent_compute
Image={{ edpm_telemetry_ceilometer_compute_image }}
Network=host
User=ceilometer
SecurityLabelType=ceilometer_polling_t
{% if edpm_telemetry_healthcheck -%}
HealthCmd=/openstack/healthcheck compute
Volume=/var/lib/openstack/healthchecks/ceilometer_agent_compute:/openstack:ro,z
{% endif -%}
LogDriver=journald
Environment=KOLLA_CONFIG_STRATEGY=COPY_ALWAYS
Environment=OS_ENDPOINT_TYPE=internal
{% for vol in _ceilometer_volumes -%}
Volume={{ vol }}
{% endfor -%}
Label=managed_by=edpm_ansible

[Service]
Restart=always
TimeoutStartSec=900
TimeoutStopSec=84

[Install]
WantedBy=multi-user.target
